# ansible/files/hampuff-sms.service
[Unit]
Description=Hampuff SMS Web Service
After=network.target
Wants=network.target

[Service]
Type=exec
User=hampuff-sms
Group=hampuff-sms
WorkingDirectory=/opt/hampuff-sms
Environment=PATH=/opt/hampuff-sms/venv/bin
Environment=FLASK_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=15015
ExecStart=/opt/hampuff-sms/venv/bin/gunicorn --bind 127.0.0.1:15015 --workers 2 --worker-class sync --timeout 30 --keep-alive 2 --max-requests 1000 --max-requests-jitter 100 wsgi:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hampuff-sms

[Install]
WantedBy=multi-user.target
